# deployment/email-monitor.service
[Unit]
Description=Inbox Bridge — Email Monitor Bot Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/inbox-bridge
Environment="PATH=%h/inbox-bridge/venv/bin"
ExecStart=%h/inbox-bridge/venv/bin/python %h/inbox-bridge/run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inbox-bridge

[Install]
WantedBy=default.target
